# rc.local
#!/bin/sh -e

nft -f /etc/firewall_rules/nft-rules.nft
exit 0
